# app-admin/le/files/le-monitor
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms fo the GNU General Public License v2
# $Header: $

LE_MONITOR_LOGFILE=${LE_MONITOR_LOGFILE:-/var/log/${SVCNAME}.log}
LE_MONITOR_PIDFILE=${LE_MONITOR_PIDFILE:-/run/${SVCNAME}.pid}
LE_MONITOR_BINARY=${LE_MONITOR_BINARY:-/usr/bin/le-monitor}
LE_MONITOR_OPTS=${LE_MONITOR_OPTS:-}

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --background \
		--exec "${LE_MONITOR_BINARY}" \
		--pidfile "${LE_MONITOR_PIDFILE}" \
		--stdout "${LE_MONITOR_LOGFILE}" \
		--stderr "${LE_MONITOR_LOGFILE}"
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop \
		--exec "${LE_MONITOR_BINARY}" \
		--pidfile "${LE_MONITOR_PIDFILE}"
	eend $?
}
