# app-backup/rackspace-driveclient/files/init.rackspace-driveclient
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or
# later
# $Header: /var/cvsroot/gentoo-x86/net-im/bitlbee/files/bitlbee.initd,v 1.3 2011/12/04 21:18:10 radhermit Exp $

pidfile="${PIDFILE:-/var/run/rackspace-driveclient.pid}"

command_args="-d"
command_args="${command_args} ${EXTRA_OPTS}"

command="/opt/rackspace-driveclient/bin/driveclient"

depend() {
  need net
  use logger
}

start() {
  ebegin "Starting ${SVCNAME}"
  start-stop-daemon --start --background --pidfile ${pidfile} ${command} -- ${command_args}
  eend $? "Failed to start ${SVCNAME}"
}

stop() {
  ebegin "Stopping ${SVCNAME}"
  start-stop-daemon --stop --pidfile ${pidfile} --name driveclient
  eend $? "Failed to stop ${SVCNAME}"
}
